# bin_rpi/crontab
* * * * * /home/pi/github-home/bin_rpi/nasMon/pingTest.sh 192.168.0.104 nas

